# docker/dev/supervisord/jupyter.conf
[program:jupyter]
directory=/home/username
command=runuser -l username -s /bin/bash -c "cd /home/username && jupyter notebook --ip=0.0.0.0 --port 8888 --no-browser"