# docker/dev/supervisord/sshd.conf
[program:sshd]
directory=/home/username
command=/usr/sbin/sshd -D